# systemd/gpstracker.service
[Unit]
Description=GPS Tracking Service through Python
After=multi-user.target

[Service]
Type=simple
User=root
Restart=always
WorkingDirectory=/home/pi/gpstracker
ExecStart=/usr/bin/python3 /home/pi/gpstracker/gps_db_save.py

[Install]
WantedBy=multi-user.target
